# daily.cron
0 5 * * * python every-day-with-jesus.py
